# contrib/DEBIAN/postinst
#!/bin/bash
#
# postinst script for p4d

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package

COMMAND="CREATE DATABASE IF NOT EXISTS p4 charset utf8;
     CREATE USER IF NOT EXISTS 'p4'@'localhost' IDENTIFIED BY 'p4';
     GRANT INDEX, CREATE ROUTINE, ALTER ROUTINE, TRIGGER, REFERENCES, LOCK TABLES, INDEX, CREATE VIEW, EXECUTE, CREATE TEMPORARY TABLES, ALTER, SELECT, INSERT, DELETE, CREATE, DROP, UPDATE ON p4.* TO 'p4'@'localhost' IDENTIFIED BY 'p4';
     GRANT INDEX, CREATE ROUTINE, ALTER ROUTINE, TRIGGER, REFERENCES, LOCK TABLES, INDEX, CREATE VIEW, EXECUTE, CREATE TEMPORARY TABLES, ALTER, SELECT, INSERT, DELETE, CREATE, DROP, UPDATE ON p4.* TO 'p4'@'%' IDENTIFIED BY 'p4';
     FLUSH PRIVILEGES;"

case "$1" in
    configure)
        echo "postinst called for configure" >&2
        if pgrep mysqld; then
           echo "Enter your database root password:" >&2
           echo ${COMMAND} | mysql -u root -Dmysql -p
        else
           echo "No local databse found, skipping db setup" >&2
        fi
        systemctl daemon-reload && systemctl enable p4d
        systemctl daemon-reload && systemctl enable w1mqtt
        systemctl daemon-reload && systemctl enable mosquitto-log.service
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
